# scripts/mypi-dreamer.service
[Unit]
Description=Run mypi dreamer cycle

[Service]
Type=oneshot
ExecStart=%h/.pi/agent/dreamer-cron.sh
